# Full Source/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24
autorestart=true
priority=1

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -rfbport 5901 -nopw -forever -shared
autorestart=true
priority=5

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ 8080 localhost:5901
autorestart=true
priority=6

[program:fluxbox]
command=/usr/bin/fluxbox
autorestart=true
environment=DISPLAY=:1
priority=8

[program:rawrxd-backend]
command=/bin/bash -lc 'if [ -f /opt/rawrxd/backend/RawrEngine.py ]; then exec python3 /opt/rawrxd/backend/RawrEngine.py; elif [ -f /opt/rawrxd/backend/chat_server.py ]; then exec python3 /opt/rawrxd/backend/chat_server.py; else echo "Backend entrypoint not found" >&2; exit 1; fi'
autorestart=true
directory=/opt/rawrxd
priority=9

[program:rawrxd-ide]
command=/usr/bin/wine /opt/rawrxd/RawrXD-Win32IDE.exe
autorestart=true
environment=DISPLAY=:1,WINEDEBUG=-all
priority=10
