# OrganizedPiProject/projects/unified-workspace/rawrz-platform/ops/true.service
# /etc/systemd/system/true.service
[Unit]
Description=true
After=network.target

[Service]
Environment=NODE_ENV=production
Environment=PORT=8080
EnvironmentFile=/opt/true/.env
WorkingDirectory=/opt/true
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
User=true
Group=true

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/true/uploads

[Install]
WantedBy=multi-user.target
